# files/hostapd/docker.timer
[Unit]
Description=Run docker 3.5 minutes after the boot

[Timer]
OnBootSec=210seconds

[Install]
WantedBy=timers.target
