# files/hostapd/hostapd-wifi.timer
[Unit]
Description=Run hostapd-wifi 3 minutes after the boot

[Timer]
OnBootSec=3min

[Install]
WantedBy=timers.target
