# files/hostapd/hostapd-wifi.timer
[Unit]
Description=Run hostapd-wifi 2 minutes after the boot

[Timer]
OnBootSec=2min

[Install]
WantedBy=timers.target
