# files/hostapd/hostapd-network.timer
[Unit]
Description=Run hostapd-network 2 minutes after the boot

[Timer]
OnBootSec=2min

[Install]
WantedBy=timers.target
