# files/hostapd/hostapd-network.timer
[Unit]
Description=Run hostapd-network 1 minutes after the boot

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
