# files/hostapd/hostapd-wifiguests.timer
[Unit]
Description=Run hostapd-wifiguests 4 minutes after the boot

[Timer]
OnBootSec=4min

[Install]
WantedBy=timers.target
